# duster.service
[Unit]
Description=Duster daemon

[Service]
Type=simple
ExecStart=(path to duster)/duster (ip:port)

[Install]
WantedBy=multi-user.target
